# vendor/xssbot/context/supervisord.conf
[supervisord]
nodaemon=true

[program:init]
command=/bin/bash -c "supervisorctl start redis && supervisorctl start marvin"
autoresetart=false

[program:marvin]
user=marvin
directory=/app
command=node src/main.js
autostart=false
autoresetart=true

[program:redis]
command=redis-server
autostart=false
autoresetart=true