# vendor/xssbot/context/supervisord.conf
[supervisord]
nodaemon=true

[program:init]
command=/bin/bash -c "iptables-legacy -t nat -A PREROUTING -p tcp --dport 80 -j REDIRECT --to-port 8000 && supervisorctl start redis && supervisorctl start $([ $NODE_ENVIRONMENT == development ] && echo 'marvin_debug' || echo 'marvin')"
autoresetart=false

[program:marvin_debug]
user=marvin
directory=/app
command=node src/main.js
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
autostart=false
autoresetart=false

[program:marvin]
user=marvin
directory=/app
command=node src/main.js
autostart=false
autoresetart=true

[program:redis]
command=redis-server
autostart=false
autoresetart=true